# root/etc/s6-overlay/s6-rc.d/svc-comfyui/run
#!/command/with-contenv sh

set -euo pipefail

# Activate virtual environment if it exists
if [ -d "/opt/venv" ]; then
    echo "[INFO] Activating virtual environment at /opt/venv"
    export VIRTUAL_ENV=/opt/venv
    export PATH="$VIRTUAL_ENV/bin:$PATH"
fi

cd /default-comfyui-bundle/ComfyUI
exec python main.py --port=8188 --listen=0.0.0.0 --cpu --enable-manager --base-directory=/workspace
