# root/etc/s6-overlay/s6-rc.d/svc-comfyui-deps/run
#!/bin/bash

set -euo pipefail

echo "[INFO] === Running preload-requirements.sh ==="

# Activate virtual environment if it exists
if [ -d "/opt/venv" ]; then
    echo "[INFO] Activating virtual environment at /opt/venv"
    export VIRTUAL_ENV=/opt/venv
    export PATH="$VIRTUAL_ENV/bin:$PATH"
fi

# Run preload-requirements.sh from custom_nodes directory
cd /custom_nodes
bash /runner-scripts/preload-requirements.sh

echo "[INFO] === preload-requirements.sh completed ==="